# roles/frontend/files/systemd-system-cerbot-renewal.timer
[Unit]
Description=Certbot Renewal Timer

[Timer]
OnBootSec=300
OnUnitActiveSec=12h

[Install]
WantedBy=multi-user.target
